# shadows_root/home/shadows/.config/systemd/user/midnight.service
[Unit]
Description=Run my script at midnight

[Service]
Type=oneshot
ExecStart=/home/shadows/.local/bin/midnight.sh
