# shadows_root/home/shadows/.config/systemd/user/midnight.timer
[Unit]
Description=Run my script at 00:00 every day

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
